# systemd/atmosbird-analysis.timer
[Unit]
Description=AtmosBird Analysis Timer (every 15 minutes)
Requires=atmosbird-analysis.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
