# systemd/database-backup.service
[Unit]
Description=EMSN PostgreSQL Database Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn2
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/maintenance/database_backup.py
StandardOutput=append:/mnt/usb/logs/database_backup.log
StandardError=append:/mnt/usb/logs/database_backup.error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
