# systemd/emsn-yearly-report.service
[Unit]
Description=EMSN Yearly Bird Report Generator
After=network.target

[Service]
Type=oneshot
User=ronny
Group=ronny
WorkingDirectory=/home/ronny/emsn2/scripts/reports
EnvironmentFile=/home/ronny/.env.emsn-reports
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/reports/yearly_report.py --pending
StandardOutput=append:/mnt/usb/logs/emsn-yearly-report.log
StandardError=append:/mnt/usb/logs/emsn-yearly-report.error.log

[Install]
WantedBy=multi-user.target
