# systemd/emsn-monthly-report.timer
[Unit]
Description=EMSN Monthly Report Generator Timer
Requires=emsn-monthly-report.service

[Timer]
# Run on the 1st of every month at 08:00
OnCalendar=*-*-01 08:00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
