# systemd/emsn-weekly-report.service
[Unit]
Description=EMSN Weekly Bird Activity Report Generator
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn2/scripts/reports
EnvironmentFile=/home/ronny/.env.emsn-reports
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/reports/weekly_report.py
StandardOutput=append:/mnt/usb/logs/emsn-weekly-report.log
StandardError=append:/mnt/usb/logs/emsn-weekly-report.error.log

[Install]
WantedBy=multi-user.target
