# systemd/emsn-reports-api.service
[Unit]
Description=EMSN Reports Web API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ronny
WorkingDirectory=/home/ronny/emsn2/reports-web
Environment="FLASK_APP=api.py"
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/reports-web/api.py
Restart=always
RestartSec=10
StandardOutput=append:/mnt/usb/logs/emsn-reports-api.log
StandardError=append:/mnt/usb/logs/emsn-reports-api.error.log

[Install]
WantedBy=multi-user.target
