# systemd/emsn-seasonal-report-winter.timer
[Unit]
Description=EMSN Winter Report Generator Timer (runs March 1st)

[Timer]
OnCalendar=*-03-01 07:00:00
Persistent=true
Unit=emsn-seasonal-report.service

[Install]
WantedBy=timers.target
