# systemd/emsn-seasonal-report-autumn.timer
[Unit]
Description=EMSN Autumn Report Generator Timer (runs December 1st)

[Timer]
OnCalendar=*-12-01 07:00:00
Persistent=true
Unit=emsn-seasonal-report.service

[Install]
WantedBy=timers.target
