# systemd/hardware-metrics.service
[Unit]
Description=EMSN Hardware Metrics Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ronny
ExecStart=/usr/bin/python3 /home/ronny/emsn2/scripts/hardware_metrics.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
