# systemd/homer-stats-update.timer
[Unit]
Description=Update Homer dashboard with EMSN stats every 5 minutes
Documentation=https://github.com/RonnyCHL/emsn2

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
