# systemd/mqtt-failover.timer
[Unit]
Description=EMSN MQTT Failover Check Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
