# systemd/nestbox-screenshot.timer
[Unit]
Description=Nestkast Screenshot Timer - elk uur

[Timer]
# Elk uur een screenshot + ML detectie
OnCalendar=*-*-* *:00:00

Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
