# systemd/network-monitor.service
[Unit]
Description=EMSN Network Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn2
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/monitoring/network_monitor.py
StandardOutput=append:/mnt/usb/logs/network_monitor.log
StandardError=append:/mnt/usb/logs/network_monitor.error.log

[Install]
WantedBy=multi-user.target
