# systemd/timer-timeline.service
[Unit]
Description=EMSN Timer Timeline Collector
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
WorkingDirectory=/home/ronny/emsn2/scripts
Environment="EMSN_DB_PASSWORD=REDACTED_DB_PASS"
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/timer_timeline.py --station zolder
StandardOutput=null
StandardError=append:/mnt/usb/logs/timer-timeline.error.log

[Install]
WantedBy=multi-user.target
