# systemd/ulanzi-bridge.service
[Unit]
Description=EMSN Ulanzi Bridge Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ronny
Group=ronny
WorkingDirectory=/home/ronny/emsn2
Environment=PYTHONUNBUFFERED=1
# Gebruik venv voor vocalization classifier dependencies (torch, librosa)
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/ulanzi/ulanzi_bridge.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ulanzi-bridge

[Install]
WantedBy=multi-user.target
