# systemd/vocalization-enricher.service
[Unit]
Description=EMSN Vocalization Enricher Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ronny
Group=ronny
WorkingDirectory=/home/ronny/emsn2
Environment=PYTHONUNBUFFERED=1
# Gebruik venv voor PyTorch/librosa dependencies
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/vocalization/vocalization_enricher.py --continuous --interval 5
Restart=always
RestartSec=60
# Memory limits om RAM problemen te voorkomen
MemoryMax=1G
MemoryHigh=800M
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vocalization-enricher

[Install]
WantedBy=multi-user.target
