# systemd/weekly-system-report.service
[Unit]
Description=EMSN Weekly System Report
Documentation=https://github.com/emsn/bird-detection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ronny
Group=ronny
WorkingDirectory=/home/ronny/emsn2/scripts/maintenance

# Execute weekly report script
ExecStart=/home/ronny/emsn2/venv/bin/python3 /home/ronny/emsn2/scripts/maintenance/weekly_system_report.py

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weekly-system-report

# Timeouts
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
